# browser-control/home/pi/.config/autostart/webcam1.desktop
[Desktop Entry]
Type=Application
Exec=curl "http://localhost/camsel.php?cam=default"
